# os-brainux/override/lib/systemd/system/rndis_gadget.service
[Unit]
Description=Enable RNDIS USB Gadget

[Service]
Type=oneshot
ExecStart=/usr/bin/enable_rndis_gadget

[Install]
WantedBy=multi-user.target
